# src/build/modules/jenkins/files/etc/supervisor/conf.d/jenkins.conf
[program:jenkins]
command = exec su jenkins -c "JENKINS_HOME=/jenkins/data java -jar /usr/share/jenkins/jenkins.war"
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
