# src/jenkins/build/modules/jenkins/files/etc/supervisor/conf.d/jenkins.conf
[program:jenkins]
command = /bin/su jenkins -c "JENKINS_HOME=/jenkins/data java -Dfile.encoding=UTF8 -Dhudson.slaves.WorkspaceList=_ -jar /usr/share/jenkins/jenkins.war"
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
